# usr/lib/systemd/system/flash-bootimage-switch-slot.service.d/99-akita.conf
[Unit]
Requires=sysinit.target
After=sysinit.target dbus.socket systemd-journald.socket 

[Service]
ExecStart=
ExecStart=/usr/share/WHO53/flash-bootimage-kernel-flash
