# debian/adaptation-lenovo-akita-configs.led.service
[Unit]
Description=LED
After=phosh.service

[Service]
ExecStart=/usr/share/WHO53/led

[Install]
WantedBy=graphical.target
