# usr/lib/systemd/system/phosh.service.d/90-akita.conf
[Service]
Environment=WLR_HWC_SYSFS_BACKLIGHT=1
